# containers/php81-cli/conf/etc/crontabs/root.example
* * * * * cd /var/www/code && php artisan schedule:run