# containers/php81-cli/conf/etc/crontabs/root.example
* * * * * php /var/www/code/artisan schedule:run